# etc/supervisor.d/docker.conf
[program:docker]
command=/usr/local/bin/wrapdocker
environment=DOCKER_DAEMON_ARGS="--insecure-registry ci1.camunda.loc:5000"
