# etc/supervisor.d/wrapdocker.conf
[program:wrapdocker]
command=/usr/local/bin/wrapdocker
environment=DOCKER_DAEMON_ARGS="--insecure-registry camunda-ci1:5000"
